# debian/postinst
#!/bin/bash
set -e

# Habilitando o serviço do maratona-firewall para ser executado no boot do
# sistema.
systemctl enable maratona-firewall.service

echo "maratona-firewall will take effect after a reboot"

exit 0
