# debian/postrm
#!/bin/bash

systemctl disable maratona-firewall.service

[ -f /lib/systemd/system/maratona-firewall.service ] && \
	rm /lib/systemd/system/maratona-firewall.service
