# debian/prerm
#!/bin/bash
set -e

systemctl disable maratona-firewall.service
