# debian/prerm
#!/bin/bash

systemctl disable maratona-firewall.service
